# morgothd/morgoth.service.in
[Unit]
Description=Morgoth Daemon

[Service]
Type=dbus
ExecStart=${CMAKE_INSTALL_PREFIX}/bin/morgothd --config /etc/morgoth.conf
Restart=always
BusName=org.morgoth
StandardOutput=inherit
StandardError=journal

[Install]
WantedBy=multi-user.target
